# Redwood-Discord-bot/redwood-bot.service
[Unit]
Description=Redwood Peak Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=redwood
WorkingDirectory=/opt/redwood-bot
EnvironmentFile=/etc/redwood-bot.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
